# autostart/Repeat-rate.desktop
[Desktop Entry]
Name=Repeat-rate
Comment=Run xset r rate 200 60
Exec=xset r rate 200 60
Terminal=false
Type=Application
