# autostart/Repeat-rate.desktop
[Desktop Entry]
Name=Repeat-rate
Comment=Run xset r rate 170 70
Exec=xset r rate 170 70
Terminal=false
Type=Application
