# autostart/Repeat-rate.desktop
[Desktop Entry]
Name=Repeat-rate
Comment=Run xset r rate 200 45
Exec=xset r rate 200 45
Terminal=false
Type=Application
