# autostart/Super.desktop
[Desktop Entry]
Name=Repeat-rate
Comment=Source .Xmodmap
Exec=xmodmap ~/.Xmodmap
Terminal=false
Type=Application
